# we-term.service
[Unit]
Description=we-term web terminal
After=network.target

[Service]
Type=simple
User=john
WorkingDirectory=/home/john/git/we-term
ExecStart=/home/john/git/we-term/venv/bin/python /home/john/git/we-term/server.py
Restart=on-failure
RestartSec=3
TimeoutStopSec=5
Environment=WETERM_HOST=0.0.0.0
Environment=WETERM_PORT=9090

[Install]
WantedBy=multi-user.target
